# judge/init/judge.conf
# Sicily Online Judge Service

description     "Sicily Online Judge"
author          "Hovey <hoveychen@gmail.com>"

start on (net-device-up
          and local-filesystems
          and started mysql
	  and runlevel [2345])
stop on runlevel [016]

respawn

env HOME=/home/judge

exec "${HOME}"/bin/judge -c "${HOME}"/judge.conf

